# backup-rotator.service
[Unit]
Description=Backup Rotator - Automated backup file rotation with Telegram notifications
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/backup-rotator
Environment="PATH=/opt/backup-rotator/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/backup-rotator/.venv/bin/python /opt/backup-rotator/run.py --config /opt/backup-rotator/config.yaml.local --execute
Restart=always
RestartSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=backup-rotator

[Install]
WantedBy=multi-user.target
